# gardener.service
[Unit]
Description=GitHub Gardener — genera un haiku e fa commit

[Service]
Type=oneshot
User=romano
WorkingDirectory=/home/romano/github-gardener
ExecStart=/usr/bin/python3 /home/romano/github-gardener/gardener.py
StandardOutput=journal
StandardError=journal
